# server/webhooks/ufactory_docs.service
[Unit]
Description=UFactory DOCS Webhooks Service
After=network.target

[Service]
ExecStart=/home/vitepress_docs/server/webhooks/start_webhooks.sh
Restart=always
User=root
Group=root

[Install]
WantedBy=multi-user.target